# scripts/.config/systemd/user/theme-switch.service
[Unit]
Description=Switch theme script


[Service]
ExecStart=/bin/bash %h/theme-switch.sh
Type=oneshot
